# deploy/maxbot.service
[Unit]
Description=MaxBot IRC Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=maxbot
WorkingDirectory=/opt/maxbot
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=on-failure
RestartSec=5
EnvironmentFile=/etc/maxbot/.env

[Install]
WantedBy=multi-user.target
